# com.medprep.api-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.medprep.api-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replaced by install_monitors.sh with actual paths -->
        <string>VENV_PYTHON</string>
        <string>-m</string>
        <string>src.monitoring.scheduler</string>
        <string>--mode</string>
        <string>api</string>
        <string>--data-dir</string>
        <string>DATA_DIR</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PROJECT_DIR</string>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run daily at 06:00 -->
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>LOG_DIR/api-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/api-monitor-error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>PROJECT_DIR</string>
    </dict>
</dict>
</plist>
